# dot_backup/coolify/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=unexpected
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hermes-webui]
command=/opt/hermes-venv/bin/python /apptoo/server.py
directory=/home/ubuntu/dev
user=ubuntu
environment=HOME="/home/ubuntu",USER="ubuntu",SHELL="/bin/zsh",PATH="/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/opt/hermes-venv/bin:/usr/local/bin:/usr/bin:/bin",HERMES_HOME="/home/ubuntu/.hermes",HERMES_WEBUI_STATE_DIR="/home/ubuntu/.hermes/webui",HERMES_WEBUI_DEFAULT_WORKSPACE="/home/ubuntu/dev",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_PYTHON="/opt/hermes-venv/bin/python"
autorestart=unexpected
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
