# dot_backup/coolify/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
priority=10

[program:hermes-webui]
command=npm start
directory=/opt/hermes-webui
user=ubuntu
environment=HOME="/home/ubuntu",USER="ubuntu",PORT="8787",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_HOST="0.0.0.0"
autorestart=true
priority=20
